# deploy/com.omoromachi.claude-allow-watcher.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.omoromachi.claude-allow-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/claude_allow_watcher</string>
    <string>--no-prompt</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>/tmp/claude_allow_watcher.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude_allow_watcher.error.log</string>
</dict>
</plist>
